# home/config/autostart/syndaemon.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/syndaemon -i .5 -k -K -R
Hidden=false
X-GNOME-Autostart-enabled=true
Name[en_US]=syndaemon
Name=syndaemon
Comment[en_US]=prevent mousing while typing
Comment=prevent mousing while typing
